# claude-dashboard-retention.timer
[Unit]
Description=Claude Dashboard — weekly retention sweep timer

[Timer]
# Weekly at 03:30 local time. Randomised by 15 min so fleet installs don't
# hammer a single NTP second when multiple users run this.
OnCalendar=Sun 03:30:00
RandomizedDelaySec=15min
Persistent=true
Unit=claude-dashboard-retention.service

[Install]
WantedBy=timers.target
