# claude-dashboard-retention.service
[Unit]
Description=Claude Dashboard — retention sweep (delete old sessions)
After=claude-dashboard.service
Requires=claude-dashboard.service

[Service]
Type=oneshot
User=hardcoremonk
# RETENTION_DAYS controls how many days of data to keep. Default 365.
# Override by dropping a systemd unit edit: `systemctl edit claude-dashboard-retention.service`
Environment=RETENTION_DAYS=365
Environment=DASHBOARD_URL=http://127.0.0.1:8765

# Uses the admin retention endpoint. If DASHBOARD_PASSWORD is set on the
# main unit, copy the same value here (or use a shared EnvironmentFile).
ExecStart=/bin/bash -c ' \
  set -eu; \
  echo "Retention sweep: deleting sessions older than $${RETENTION_DAYS}d"; \
  curl -fsSL -X DELETE \
    "$${DASHBOARD_URL}/api/admin/retention?older_than_days=$${RETENTION_DAYS}&confirm=true" \
  || { echo "retention request failed"; exit 1; }; \
  echo "Incremental vacuum reclaim …"; \
  sqlite3 /home/hardcoremonk/.claude/dashboard.db "PRAGMA incremental_vacuum;" || true; \
'

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/hardcoremonk/.claude
