# claude-dashboard.service
[Unit]
Description=Claude Usage Dashboard
After=network.target

[Service]
Type=simple
# User/paths are hardcoded to the install location — edit these two lines
# when relocating or running as a different user.
User=hardcoremonk
WorkingDirectory=/home/hardcoremonk/projects/claude-dashboard
ExecStart=/home/hardcoremonk/projects/claude-dashboard/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8765 --loop asyncio --http h11 --log-level info
Restart=always
RestartSec=5
Environment=PORT=8765
# Optional: set DASHBOARD_PASSWORD=yourpassword for HTTP Basic Auth

# Resource limits — tuned for a 500 MB-ish DB and 8 parallel scan workers.
# Raise MemoryMax / TasksMax if the watcher logs OOMs.
MemoryMax=512M
MemoryHigh=384M
CPUQuota=150%
TasksMax=256
LimitNOFILE=8192

# Hardening — read-only filesystem except for the DB dir and bytecode cache.
# .venv is NOT writable at runtime: start.sh installs deps before handing off
# to uvicorn, so the service itself never needs pip. If you change
# requirements.txt, restart the service (start.sh will reinstall).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/hardcoremonk/.claude /home/hardcoremonk/projects/claude-dashboard/__pycache__
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-dashboard

[Install]
WantedBy=multi-user.target
